# templates/kafka-broker.service.j2
[Unit]
Description=Kafka broker
After=network.target

[Service]
User=kafka
Group=kafka

WorkingDirectory={{ kafka_home_dir }}

ExecStart={{ kafka_home_dir }}/bin/kafka-server-start.sh {{ kafka_broker_config_path }}

Restart=on-failure
RestartSec=10

SuccessExitStatus=143

[Install]
WantedBy=multi-user.target